# misc/cm_daily.timer
[Unit]
Description=Daily tasks for ChemManager

[Timer]
OnCalendar=3:30
Unit=cm_daily.service

[Install]
WantedBy=basic.target
