# chemman/core/templates/_deploy/chemman.socket
[Unit]
Description=Gunicorn socket for ChemManager

[Socket]
ListenStream={{ socket }}
SocketUser={{ http_user }}

[Install]
WantedBy=sockets.target
